# monstercatyoutubenotifier.service
[Unit]
Description=MonstercatYoutubeNotifier

[Service]
Type=simple
ExecStart=/root/git/MonstercatYoutubeNotifier/notifier.py cron
